# recipes-core/hello/files/hello.service
[Unit]
Description=Hello World Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Hello from Yocto!"

[Install]
WantedBy=multi-user.target
